# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:gunicorn]
command=/bin/sh -c "gunicorn -w 5 -b 0.0.0.0:8080 --access-logfile - --error-logfile - --log-level WARNING --timeout 600 app.app:app"
directory=/usr/src/app
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:celery_worker]
command=/bin/sh -c "celery -A app.celery_app.celery worker --loglevel=WARNING"
directory=/usr/src/app
broker_connection_retry_on_startup=true
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:celery_beat]
command=/bin/sh -c "celery -A app.celery_app.celery beat --loglevel=WARNING --pidfile=/tmp/celerybeat.pid --schedule=/tmp/celerybeat-schedule"
directory=/usr/src/app
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
